# Linux/deployment/ansible-wazuh-agent/files/osquery.conf
{
    "options": {
        "config_plugin": "filesystem",
        "logger_plugin": "filesystem",
        "logger_path": "/var/log/osquery",
        "disable_logging": "false",
        "pidfile": "/var/osquery/osquery.pidfile",
        "disable_events": "false",
        "disable_audit": "false",
        "audit_allow_config": "true",
    	"audit_persist": "true",
    	"audit_allow_process_events": "true",
    	"audit_allow_fim_events": "true",
    	"audit_allow_user_events": "true",
    	"audit_allow_sockets": "true",
        "schedule_splay_percent": "10",
        "database_path": "/var/osquery/osquery.db",
        "utc": "true",
        "worker_threads": "3"
    },
    "yara": {
        "signatures": {
        "sig_group_1": [ "yara.sig" ]
        },
       "file_paths": {
       "homes": [ "sig_group_1" ],
       "systembins": [ "sig_group_1" ],
       "tmp": [ "sig_group_1" ],
       "etc": [ "sig_group_1" ],
       "www": [ "sig_group_1" ]
        }
     },

    "file_paths": {
        "homes": [ "/home/%%" ],
        "etc": [ "/etc/%%" ],
        "system_binaries": [ "/usr/bin/%", "/usr/sbin/%", "/usr/local/sbin/%", "/usr/local/bin/%", "/sbin/%", "/bin/%", "/usr/games/%", "/usr/local/games/%" ],
        "tmp": [ "/tmp/%%" ],
        "www": [ "/var/www/%%" ]
	
    },

  "schedule": {
    "yara": {
      "query": "SELECT * FROM yara_events",
      "interval": 30,
      "removed": false
     },
    "listening_ports": {
      "query": "SELECT DISTINCT process.name, process.path, listening.port, listening.address, process.pid FROM processes AS process JOIN listening_ports AS listening ON process.pid = listening.pid where port!=0",
      "interval": 60,
      "removed": true
    },
    "outbound_comms": {
      "query": "select s.pid, p.name, local_address, remote_address, family, protocol, local_port, remote_port from process_open_sockets s join processes p on s.pid = p.pid where family = 2",
      "interval": 60,
      "removed": false
    },
    "deleted_but_running": {
      "query": "SELECT name, path, pid FROM processes WHERE on_disk = 0",
      "interval": 300,
      "removed": true
    },
    "top_processes": {
      "query": "select count(pid) as total, name from processes group by name order by total desc limit 10",
      "interval": 1800,
      "removed": false
    },

    "new_listening_port_process": {
      "query": "select distinct process.name, listening.port, listening.address, process.pid from processes as process join listening_ports as listening on process.pid = listening.pid",
      "interval": 300
    },
    "suspicious_outbound": {
      "query": "select s.pid, p.name, local_address, remote_address, family, protocol, local_port, remote_port from process_open_sockets s join processes p on s.pid=p.pid where remote_port not in (53,80,443, 1514) and family=2",
      "interval": 600
    },

    "kmodules": {
      "query": "select name from kernel_modules",
      "interval": 120,
      "removed": false
    },
    "deb_packages": {
      "query": "SELECT name, version, source FROM deb_packages",
      "interval": 120,
      "removed": true
    },
    "rpm_packages": {
      "query": "SELECT name, version, source FROM rpm_packages",
      "interval": 120,
      "removed": true
    },
    "iptables_monitor": {
      "query": "select chain,protocol,target from iptables",
      "interval": 60,
      "removed": true
    },
    "crontab_monitor": {
      "query": "select command, path from crontab",
      "interval": 60,
      "removed": false
    },
    "etc_monitor": {
      "query": "SELECT path, hash.sha256 FROM file JOIN hash USING (path) WHERE file.directory = '/etc/'",
      "interval": 60,
      "removed": true
    },
    "suid_monitor": {
      "query": "select * from suid_bin",
      "interval": 60,
      "removed": true
    },
    "login_monitor": {
      "query": "SELECT username, p.name AS process, DATETIME(time, 'unixepoch') AS datetime, host FROM last AS l LEFT JOIN processes AS p ON p.pid = l.pid WHERE host <> '' ORDER BY time DESC",
      "interval": 60,
      "removed": false
    },
    "user_monitor": {
      "query": "select username,shell from users",
      "interval": 60,
      "removed": true
    },
    "history_monitor": {
      "query": "select history_file,command from users join shell_history using (uid)",
      "interval": 60,
      "removed": true
    },
    "authorized_keys_monitor": {
      "query": "select key_file from users join authorized_Keys using (uid)",
      "interval": 60,
      "removed": false
    },
    "sudoers_monitor": {
      "query": "select * from sudoers",
      "interval": 60
    },
    "startup_items": {
      "query": "select name,path from startup_items where status='enabled'",
      "interval": 60
    },
    "apparmor_monitor": {
      "query": "select path,mode from apparmor_profiles",
      "interval": 300,
      "removed": false
    },
    "aptsources_monitor": {
      "query": "select source,base_uri from apt_sources",
      "interval": 60,
      "removed": false
    },
    "arpcache_monitor": {
      "query": "select address from arp_cache",
      "interval": 60,
      "removed": false
    },
    "sshkey_monitor": {
      "query": "select username,path from users join user_ssh_keys using (uid)",
      "interval": 600,
      "removed": false
    },
    "dns_monitor": {
      "query": "select * FROM dns_resolvers",
      "interval": 120,
      "removed": false
    },
    "user_logins": {
      "query": "select * from logged_in_users where type!='dead'",
      "interval": 60,
      "removed": false
    },
    "python_packages": {
      "query": "select name from python_packages",
      "interval": 60,
      "removed": false
    },
    "password_changes": {
      "query": "select password_status,username,last_change from shadow where password_status!='locked'",
      "interval": 60,
      "removed": false
    },

    "system_info": {
    "query": "SELECT hostname, cpu_brand, physical_memory FROM system_info",
    "interval": 3600
    },

    "high_load_average": {
    "query": "SELECT period, average, '70%' AS 'threshold' FROM load_average WHERE period = '15m' AND average > '0.7'",
    "interval": 900,
    "description": "Report if load charge is over 70 percent."
    }
  },

    "packs": {
   	"osquery-monitoring": "/opt/osquery/share/osquery/packs/osquery-monitoring.conf",
        "incident-response": "/opt/osquery/share/osquery/packs/incident-response.conf",
        "vuln-management": "/opt/osquery/share/osquery/packs/vuln-management.conf",
        "ossec-rootkit": "/opt/osquery/share/osquery/packs/ossec-rootkit.conf"
    }
}
