# Linux/deployment/wazuh-manager/files/local_rules.xml
<!--
  - custom & modified wazuh rules for CCDC usage
  - author: jbokor@uci.edu, UCI CCDC
  - 02/2022
--> 
<group name="CUSTOM RULES">


<!-- ======================================== -->
<!-- MODIFIED RULES -->
<!-- ======================================== -->
    <rule id="5303" level="5" overwrite="yes">
      <if_sid>5300</if_sid>
      <regex>session opened for user root|^'su root'|</regex>
      <regex>^+ \S+ \S+\proot$|^\S+ to root on|^SU \S+ \S+ + \S+ \S+-root$</regex>
      <description>User successfully changed UID to root.</description>
      <mitre>
        <id>T1078</id>
      </mitre>
      <group>authentication_success,pci_dss_10.2.5,gpg13_7.6,gpg13_7.8,gpg13_7.9,gdpr_IV_35.7.d,gdpr_IV_32.2,hipaa_164.312.b,nist_800_53_AU.14,nist_800_53_AC.7,tsc_CC6.8,tsc_CC7.2,tsc_CC7.3,</group>
    </rule>

    <rule id="10100" level="6" overwrite="yes">
      <if_group>authentication_success</if_group>
      <if_fts />
      <group>authentication_success,</group>
      <description>First time user logged in.</description>
      <mitre>
        <id>T1078</id>
      </mitre>
    </rule>


<!-- ======================================== -->
<!-- CUSTOM RULES -->
<!-- ======================================== -->
<!-- FIM -->

    <rule id="100009" level="0">
      <if_group>syscheck</if_group>
      <field name="file">/root/.viminfo</field>
      <match>modified</match>
      <description>viminfo (vim history file) updated</description>
    </rule>

    <rule id="100010" level="0">
      <if_group>syscheck</if_group>
      <match>/etc/resolv.conf</match>
      <description>resolv.conf modified, big whoop</description>
    </rule>

    <rule id="100011" level="10">
      <if_group>syscheck</if_group>
      <match>modified|added|deleted</match>
      <field name="file">^/etc/cron.</field>
      <description>File changed inside Cron!</description>
    </rule>


    <rule id="100012" level="10">
      <if_group>syscheck</if_group>
      <field name="file">/etc/pam</field>
      <description>PAM configuration modified!</description>
    </rule>

    <rule id="100013" level="9">
      <if_group>syscheck</if_group>
      <field name="file">authorized_keys$</field>
      <match>modified</match>
      <description>authorized_keys updated, UHOH</description>
    </rule>

<!-- read internal docs for more info if this isn't working as expected -->
    <rule id="100014" level="12">
      <if_group>syscheck</if_group>
      <field name="changed_content">thug|THUG</field>
      <match>modified</match>
      <description>WORD "thug" HAS BEEN ADDED TO FILE</description>
    </rule>

    <rule id="100015" level="6">
      <if_group>syscheck</if_group>
      <field name="file">^/bin | ^/usr/bin</field>
      <description>User-level binary file changed!</description>
    </rule>

    <rule id="100016" level="8">
      <if_group>syscheck</if_group>
      <field name="file">^/sbin | ^/usr/sbin</field>
      <description>Superuser binary file changed!</description>
    </rule>

    <rule id="100017" level="6">
      <if_group>syscheck</if_group>
      <field name="file">/etc/apt</field>
      <description>Apt configuration files updated</description>
    </rule>

<!-- OSQUERY -->
  <rule id="100901" level="6"> 
    <if_sid>24010</if_sid>
    <field name="osquery.name">aptsources_monitor</field>
    <description>osquery: change in apt sources, in $(osquery.columns.source)</description>
  </rule>

  <rule id="100902" level="7"> 
    <if_sid>24010</if_sid>
    <field name="osquery.name">deleted_but_running</field>
    <description>osquery: DELETED BUT RUNNING: $(osquery.columns.path) with pid $(osquery.columns.pid)</description>
  </rule>

  <rule id="100903" level="6"> 
    <if_sid>24010</if_sid>
    <field name="osquery.name">suspicious_outbound</field>
    <description>osquery: Process $(osquery.columns.p.name) has initiated outbound to remote port $(osquery.columns.remote_port)</description>
  </rule>
      
  <rule id="100904" level="7"> 
    <if_sid>24010</if_sid>
    <field name="osquery.name">sudoers_monitor</field>
    <description>osquery: sudoers file rule change</description>
  </rule>

  <rule id="100905" level="6"> 
    <if_sid>24010</if_sid>
    <field name="osquery.name">startup_items</field>
    <description>osquery: startup items change: $(osquery.columns.path)</description>
  </rule>

  <rule id="100906" level="5"> 
    <if_sid>24010</if_sid>
    <field name="osquery.name">kmodules</field>
    <description>osquery: kernel modules updated: $(osquery.columns.name)</description>
  </rule>

<!-- REMOVING ANNOYING ALERTS -->
    <rule id="101000" level="0"> 
      <if_sid>24001</if_sid> 
      <match>already running</match>
      <description>osquery already running</description> 
    </rule>


<!-- AUDITD rules -->
    <rule id="100200" level="5">
        <if_sid>80792</if_sid>
        <list field="audit.command" lookup="match_key">etc/lists/suspicious-programs</list>
        <description>Audit: Suspicious Command: $(audit.exe)</description>
        <group>audit_command,</group>
    </rule>

    <rule id="100210" level="7">
        <if_sid>80792</if_sid>
        <list field="audit.command" lookup="match_key" check_value="orange">etc/lists/suspicious-programs</list>
        <description>Audit: Suspicious Command: $(audit.exe)</description>
        <group>audit_command,</group>
    </rule>

    <rule id="100220" level="10">
        <if_sid>80792</if_sid>
        <list field="audit.command" lookup="match_key_value" check_value="red">etc/lists/suspicious-programs</list>
        <description>Audit: Highly Suspicious Command executed: $(audit.exe)</description>
        <group>audit_command,</group>
    </rule>



</group>
